# configs/systemd/user/dunst.service
[Unit]
Description=Dunst Notification Daemon
After=graphical-session.target
ExecCondition=/usr/lib/systemd/systemd-xdg-autostart-condition "wlroots:sway:Wayfire:labwc:Hyprland" ""

[Service]
Type=simple
ExecStart=/usr/bin/dunst
Restart=on-failure
Slice=app-graphical.slice

[Install]
WantedBy=graphical-session.target
