# configs/systemd/user/swww-daemon.service
[Unit]
Description=Swww Daemon
After=graphical-session.target
ExecCondition=/usr/lib/systemd/systemd-xdg-autostart-condition "wlroots:sway:Wayfire:labwc:Hyprland" ""

[Service]
Type=simple
ExecStart=/usr/bin/swww-daemon --format xrgb
Restart=on-failure
Slice=background-graphical.slice

[Install]
WantedBy=graphical-session.target

